# systemd/minifw-ai-web.service
[Unit]
Description=MiniFW-AI Web Admin Panel
After=network.target minifw-ai.service
Wants=minifw-ai.service

[Service]
Type=simple
User=root
WorkingDirectory=/opt/minifw_ai
EnvironmentFile=/etc/minifw/minifw.env
Environment=PYTHONPATH=/opt/minifw_ai/app
Environment=MINIFW_AUDIT_LOG=/opt/minifw_ai/logs/audit.jsonl
ExecStart=/opt/minifw_ai/venv/bin/uvicorn web.app:app --host 127.0.0.1 --port 8443 --ssl-keyfile /etc/minifw/tls/server.key --ssl-certfile /etc/minifw/tls/server.crt
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
